# SPECS/livepatching-subscription/livepatching.service
[Unit]
Description=Subscribe to kernel livepatching
Before=kpatch.service

[Service]
Type=simple
ExecStart=/lib/livepatching/livepatching.systemd
