# SPECS-EXTENDED/authd/auth.socket
[Unit]
Description=Authd Activation Socket

[Socket]
ListenStream=113
Accept=true

[Install]
WantedBy=sockets.target
