# SPECS/erlang/epmd.socket
[Unit]
Description=Erlang Port Mapper Daemon Activation Socket

[Socket]
ListenStream=4369
Accept=false

[Install]
WantedBy=sockets.target

