# SPECS/erlang/epmd@.socket
[Unit]
Description=Erlang Port Mapper Daemon Activation Socket

[Socket]
ListenStream=%I
Accept=false

[Install]
WantedBy=sockets.target

