# toolkit/imageconfigs/files/linuxguard/sshd-keygen.service
[Unit]
Description=Generate sshd host keys
ConditionPathExists=|!/rw/etc/ssh/ssh_host_rsa_key
ConditionPathExists=|!/rw/etc/ssh/ssh_host_ecdsa_key
ConditionPathExists=|!/rw/etc/ssh/ssh_host_ed25519_key
Before=sshd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ssh-keygen -A -f /rw

[Install]
WantedBy=multi-user.target
