# SPECS-EXTENDED/mdadm/raid-check.service
[Unit]
Description=RAID setup health check

[Service]
Type=oneshot
ExecStart=/usr/sbin/raid-check
