# base/comps/azurelinux-release/proc-version-override.service
[Unit]
Description=Override /proc/version for legacy OS detection (CBL-Mariner compat)
Before=waagent.service cloud-init.service cloud-init-local.service
After=local-fs.target
ConditionVirtualization=vm

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/proc-version-override
ExecStop=-/bin/umount /proc/version
ExecStopPost=-/bin/rm -f /run/proc_version_override

[Install]
WantedBy=multi-user.target
