# SPECS/vsftpd/vsftpd.target
[Unit]
Description=FTP daemon
After=network-online.target

[Install]
WantedBy=multi-user.target
