# SPECS/iptables/iptables.service
[Unit]
Description=Load Iptables Rules
ConditionFileIsExecutable=/etc/systemd/scripts/iptables
After=network.target

[Service]
Type=forking
ExecStart=/etc/systemd/scripts/iptables
ExecStop=/etc/systemd/scripts/iptables.stop
TimeoutSec=0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
