# specs/c/cvs/cvs.socket
[Unit]
Description=CVS Server Activation Socket
PartOf=cvs.target

[Socket]
ListenStream=2401
Accept=true

[Install]
WantedBy=sockets.target
