# SPECS/chrony/chrony-dnssrv@.timer
[Unit]
Description=Periodic DNS SRV lookup of %I for chrony

[Timer]
OnActiveSec=0
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target
