# SPECS-EXTENDED/bluez/btattach-bcm@.service
[Unit]
Description=btattach for Broadcom devices

[Service]
Type=simple
ExecStart=/usr/libexec/bluetooth/btattach-bcm-service.sh %I
