# SPECS/rear/rear.service
[Unit]
Description=Update ReaR rescue image

[Service]
Type=oneshot
ExecStart=/usr/sbin/rear checklayout || /usr/sbin/rear mkrescue
