# SPECS-EXTENDED/rp-pppoe/pppoe-server.service
[Unit]
Description=PPPoE Server.
After=syslog.target

[Service]
ExecStart=/sbin/pppoe-server

[Install]
WantedBy=multi-user.target
