# SPECS-EXTENDED/cgdcbxd/cgdcbxd.service
[Unit]
Description=cgdcbxd daemon
After=syslog.target

[Service]
ExecStart=/usr/sbin/cgdcbxd -n

[Install]
WantedBy=multi-user.target

