# SPECS/finger/finger.socket
[Unit]
Description=Finger Socket

[Socket]
ListenStream=79
Accept=yes

[Install]
WantedBy=sockets.target
