# SPECS/finger/finger@.service
[Unit]
Description=Finger Per-Connection Daemon

[Service]
ExecStart=-/usr/bin/in.fingerd
StandardInput=socket
StandardError=syslog
